# ansible/roles/zfs/systemd-scrub/files/zfs-scrub.timer
[Unit]
Description=Scrub ZFS pool weekly

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
